# contrib/seeder@.service
# Optional systemd template unit for the multi-chain DNS seeder.
# The instance name after '@' is the chain preset (chains/<name>.conf).
#
# Install + enable one chain per instance:
#   sudo cp contrib/seeder@.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now seeder@wojakcoin
#   sudo systemctl enable --now seeder@pepecoin
#   journalctl -u seeder@pepecoin -f
#
# Adjust WorkingDirectory/ExecStart paths if you cloned somewhere other than /root.
[Unit]
Description=Multi-chain DNS seeder (%i)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/wojakcoin-seeder
ExecStart=/root/wojakcoin-seeder/seederctl.sh %i run
Restart=on-failure
RestartSec=10
# Allow binding port 53 without full root (optional; comment out if undesired).
AmbientCapabilities=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
