# contrib/wojakcoin-seeder.service
# Optional systemd unit for the WojakCoin DNS seeder (single server).
#
# Install:
#   sudo cp contrib/wojakcoin-seeder.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now wojakcoin-seeder
#   journalctl -u wojakcoin-seeder -f      # follow logs
#
# Adjust WorkingDirectory/ExecStart paths if you cloned somewhere other than /root.
# Uses 'seederctl.sh run' (foreground) so systemd supervises the process directly.
[Unit]
Description=WojakCoin DNS seeder / network crawler
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/wojakcoin-seeder
ExecStart=/root/wojakcoin-seeder/seederctl.sh run
Restart=on-failure
RestartSec=10
# Allow binding port 53 without running as full root (optional; comment out if undesired).
AmbientCapabilities=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
